# launchd/com.kaimemories.onionshell.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.kaimemories.onionshell</string>

    <key>ProgramArguments</key>
    <array>
        <string>__PYTHON__</string>
        <string>__SCRIPT__</string>
        <string>_daemon</string>
    </array>

    <!-- Restart automatically if the process exits for any reason -->
    <key>KeepAlive</key>
    <true/>

    <!-- Background: prevents App Nap throttling while still being energy-efficient -->
    <key>ProcessType</key>
    <string>Background</string>

    <!-- Start at login -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Throttle rapid restart loops (e.g. if daemon crashes on startup) -->
    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>__HOME__/.onion_shell/watcher.log</string>
    <key>StandardErrorPath</key>
    <string>__HOME__/.onion_shell/watcher.log</string>
</dict>
</plist>
